# phantom.ini
[program:phantom1]
command=phantomjs /app/src/phantom.js 18080
user =root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/phantom.log

[program:phantom2]
command=phantomjs /app/src/phantom.js 18081
user =root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/phantom.log

[program:phantom3]
command=phantomjs /app/src/phantom.js 18082
user =root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/phantom.log
